# discontinued/flymon/flymon.service
[Unit]
Description=Flymon - Jeenode based monitoring

[Service]
ExecStart=/usr/bin/flymon.py

[Install]
WantedBy=multi-user.target

